# pkg/deb/debian/postinst
#!/bin/bash

set -e

source /usr/lib/growbies/paths.env

# Install into a python virtual environment
${PATH_DEBIAN_BASE_PYTHON} -m venv "${PATH_USR_LIB_GROWBIES_VENV}"
source "${PATH_USR_LIB_GROWBIES_VENV_ACTIVATE}"
pip install "${PATH_USR_LIB_GROWBIES}"/*.whl --force-reinstall
rm "${PATH_USR_LIB_GROWBIES}"/*.whl

# Create a system user 'growbies' if it doesn't exist
if ! id -u growbies >/dev/null 2>&1; then
    adduser --system --group --no-create-home --quiet growbies
fi

systemctl enable growbies
if systemd-notify --booted; then
  systemctl restart growbies
fi

